# config/systemd/user/maintenance-reminder.service
[Unit]
Description=Enable reshift

[Service]
Type=oneshot
ExecStart=/home/ashish/.scripts/systemd/maintenance.sh
