# config/systemd/user/redshift-evening.service
[Unit]
Description=Enable reshift

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=redshift -PO 4500
