# config/systemd/user/redshift-evening.timer
[Unit]
Description=Enable redshift at 22:00 daily

[Timer]
OnCalendar=*-*-* 22:00:00
Persistent=false

[Install]
WantedBy=timers.target
