# config/systemd/user/redshift-morning.service
[Unit]
Description=Disable reshift

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=redshift -x
