# config/systemd/user/trash-empty-weekly.timer
[Unit]
Description=Empty trash weekly

[Timer]
OnCalendar=Sat 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
