# config/systemd/user/trash-empty-weekly.service
[Unit]
Description=Empty trash older than 7 days

[Service]
Type=oneshot
ExecStart=gtrash prune --day 7 --force
