# config/root-NetworkManager/dispatcher.d/02-mail.sh
#!/bin/sh
case "$2" in
    up|down)
        sigdsblocks 3
        systemctl --no-block restart goimapnotify.service
        ;;
esac
